# etc/init/intel-undervolt-daemon.conf
start on started system-services
stop on stopping system-services
respawn

exec /usr/local/bin/intel-undervolt daemon
