# etc/init/intel-undervolt.conf
start on started system-services
task

exec /usr/local/bin/intel-undervolt apply
